# tools/vagrant-boxes/go_build_mpc_vm
#!/bin/bash


#
# Vagrant 2.xx installed ?
printf "=> Verification : Vagrant installé : "
if [[ $(which vagrant) ]]; then
    echo "Ok"
    printf "$(vagrant --version)"
  else
    printf "${RED}Ko${NC}\n"
    echo "Pour installer vagrant rapidement:"
    echo "wget https://releases.hashicorp.com/vagrant/2.2.5/vagrant_2.2.5_x86_64.deb"
    echo "apt install ./vagrant_2.2.5_x86_64.deb"
    exit;
fi

#
# Start VM
vagrant up

#
# Install .deb
#
init_script='
sudo apt-get update;
sudo apt-get install -y python-minimal;
sudo apt-get install -y wget net-tools;
sudo apt-get install -y docker-compose;
sudo apt install -y zip;
sudo gpasswd -a $USER docker'

vagrant ssh -c "$init_script"


#
# Reboot VM
#
vagrant halt
vagrant up

#
# /etc/init.d/
content='
! /bin/sh

### BEGIN INIT INFO
# Provides:		sshd
# Required-Start:	$remote_fs $syslog
# Required-Stop:	$remote_fs $syslog
# Default-Start:	2 3 4 5
# Default-Stop:		
# Short-Description:	OpenBSD Secure Shell server
### END INIT INFO

set -e
case "$1" in
  start)
	cd /home/vagrant/mon_premier_ctf; 
    ./go_first_install_webserver_run -y
	;;
  stop)
	;;
  *)
	exit 1
esac

exit 0
'
 
#
# Install mon_premier_ctf and run it
#
init_script='
git clone https://github.com/monpremierctf/mon_premier_ctf.git;
cd mon_premier_ctf;
./go_first_install_webserver_run -y
echo "$content" | sudo tee -a /etc/init.d/mon_premier_ctf.sh
sudo chmod 755 /etc/init.d/mon_premier_ctf.sh
sudo update-rc.d mon_premier_ctf.sh defaults'

vagrant ssh -c "$init_script"

#
# Set run at VM startup




#
# Upload VM
#
# ls -al ~/VirtualBox\ VMs/vagrant-boxes_default_1566830539113_74963/